# supervisord.conf
[supervisord]
nodaemon=true

[program:go-app]
command=go run /go/src/survey_backend/main.go
autorestart=true
