# supervisord.conf
[supervisord]
nodaemon=true

[program:go-app]
command=/usr/local/go/bin/go run /go/src/survey_backend/main.go
autorestart=true

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
